# scripts/launchd/com.gradproject.chatlog.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.gradproject.chatlog.cleanup</string>
    <key>ProgramArguments</key>
    <array>
      <string>/ABSOLUTE/PATH/TO/GRADUATION_PROJECT/scripts/clear_chatlog.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key>
      <integer>0</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/GRADUATION_PROJECT/logs/clear_chatlog.out</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/GRADUATION_PROJECT/logs/clear_chatlog.err</string>
  </dict>
</plist>
